# mysql.conf
[program:mysql]
command=/usr/sbin/mysqld --basedir=/var/lib/mysql --datadir=/var/lib/mysql
process_name=%(program_name)s
priority=1001
autostart=true
startretries=3
autorestart=true
user=root